# src/IskanderHearth/firmware/solidarity_hat/systemd/hearth-sovereignty.service
# hearth-sovereignty.service
# Solidarity HAT Hardware Kill Switch Monitor
#
# Monitors physical DPST switch states via GPIO and enforces OS-level consequences.
# On GPU kill: freezes Ollama (SIGSTOP) + halts LangGraph orchestrator.
# On hardware init failure: writes sovereignty lock, blocking Web3 signing.
#
# Install:
#   sudo cp hearth-sovereignty.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now hearth-sovereignty.service

[Unit]
Description=Iskander Hearth — Hardware Sovereignty Daemon (Kill Switch Monitor)
Documentation=file:///opt/iskander/firmware/solidarity_hat/hardware_sovereignty_daemon.py
After=network.target
Wants=network.target

# Must start before Web3 signing daemon
Before=hearth-web3-signer.service

[Service]
Type=notify
NotifyAccess=main

User=hearth-hw
Group=hearth-hw

# Needs access to GPIO (gpio group) and I2C (i2c group)
SupplementaryGroups=gpio i2c

ExecStart=/opt/iskander/venv/bin/python3 \
    /opt/iskander/firmware/solidarity_hat/hardware_sovereignty_daemon.py

Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# If daemon fails 3 times in 60s, do not restart — a human must investigate
# Sovereignty lock will have been written; Web3 signing is already suspended.

# Environment
Environment=HEARTH_GPIO_MIC_KILL=4
Environment=HEARTH_GPIO_WIFI_KILL=17
Environment=HEARTH_GPIO_GPU_KILL=27
Environment=HEARTH_KILL_POLL_SEC=0.25
Environment=HEARTH_OLLAMA_URL=http://127.0.0.1:11434
Environment=HEARTH_LANGGRAPH_URL=http://127.0.0.1:8123
Environment=HEARTH_SOVEREIGNTY_LOCK=/run/iskander/sovereignty.lock

# Runtime directory for sovereignty lock
RuntimeDirectory=iskander
RuntimeDirectoryMode=0750

# Hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/run/iskander
CapabilityBoundingSet=

# systemd watchdog: if daemon stops notifying within 30s, restart it
WatchdogSec=30

[Install]
WantedBy=multi-user.target
