# src/IskanderHearth/firmware/solidarity_hat/systemd/hearth-permacomputing.service
# hearth-permacomputing.service
# Solidarity HAT INA3221 Power Monitor & Graceful Degradation Daemon
#
# Reads solar/battery voltage from INA3221 via I2C.
# On low voltage: pauses heavy agents, broadcasts status to federated nodes.
# On critical voltage: reduces CPU performance, pauses all background agents.
# On shutdown voltage: initiates graceful OS poweroff.
#
# Install:
#   sudo cp hearth-permacomputing.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now hearth-permacomputing.service

[Unit]
Description=Iskander Hearth — Permacomputing Power Monitor (INA3221)
Documentation=file:///opt/iskander/firmware/solidarity_hat/permacomputing_monitor.py
After=network.target hearth-sovereignty.service
Wants=network.target

[Service]
Type=notify
NotifyAccess=main

User=hearth-hw
Group=hearth-hw

SupplementaryGroups=i2c

ExecStart=/opt/iskander/venv/bin/python3 \
    /opt/iskander/firmware/solidarity_hat/permacomputing_monitor.py

Restart=on-failure
RestartSec=10
StartLimitIntervalSec=120
StartLimitBurst=3

# Environment
Environment=HEARTH_I2C_BUS=1
Environment=HEARTH_POWER_POLL_SEC=30
Environment=HEARTH_BATT_WARN_V=11.5
Environment=HEARTH_BATT_CRITICAL_V=10.8
Environment=HEARTH_BATT_SHUTDOWN_V=10.2
Environment=HEARTH_ACTIVITYPUB_WEBHOOK=http://127.0.0.1:8080/internal/power-status
Environment=HEARTH_COGS_URL=http://127.0.0.1:8124/agents/cogs
Environment=HEARTH_LANGGRAPH_URL=http://127.0.0.1:8123
Environment=HEARTH_SOVEREIGNTY_LOCK=/run/iskander/sovereignty.lock

RuntimeDirectory=iskander
RuntimeDirectoryMode=0750

# Allow shutdown initiation (shutdown command requires elevated access)
# Grant only the specific capability needed
AmbientCapabilities=CAP_SYS_BOOT

# Hardening (relaxed only for CAP_SYS_BOOT)
NoNewPrivileges=no
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/run/iskander /etc/iskander
CapabilityBoundingSet=CAP_SYS_BOOT

WatchdogSec=120

[Install]
WantedBy=multi-user.target
