# src/IskanderHearth/software/thermal/cpu_thermal_manager.service
# cpu_thermal_manager.service
# Iskander Hearth — CPU Thermal Tuning via tuned
#
# Applies adaptive CPU thermal profiles using tuned(8).
# Activates the 'balanced' profile by default; on inference load,
# the gpu_thermal_manager.sh script may switch to 'throughput-performance'.
# On LOW_POWER state from permacomputing_monitor, this unit's ConditionPathExists
# gate is cleared and it falls back to 'powersave' via ExecReload.
#
# Install:
#   sudo cp cpu_thermal_manager.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now cpu_thermal_manager.service
#
# Dependencies:
#   tuned (install: dnf install tuned / apt install tuned)
#   cpupower (install: apt install linux-cpupower / dnf install kernel-tools)

[Unit]
Description=Iskander Hearth — CPU Adaptive Thermal Profile Manager
Documentation=man:tuned(8) man:tuned-adm(8)
After=tuned.service multi-user.target hearth-permacomputing.service
Requires=tuned.service

[Service]
Type=oneshot
RemainAfterExit=yes

# Apply balanced profile at boot — good default for desktop AI node
# CPU boost enabled, frequency scaling adaptive, fan curves not pinned
ExecStart=/usr/sbin/tuned-adm profile balanced

# On reload: switch to powersave (triggered by permacomputing_monitor on LOW_POWER)
ExecReload=/usr/sbin/tuned-adm profile powersave
ExecReload=/usr/bin/cpupower frequency-set --governor powersave

# On stop: restore to balanced profile and default governor
ExecStop=/usr/sbin/tuned-adm profile balanced
ExecStop=/usr/bin/cpupower frequency-set --governor ondemand

# Log profile changes to journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hearth.cpu_thermal

# Environment overrides
Environment=TUNED_DEFAULT_PROFILE=balanced
Environment=TUNED_INFERENCE_PROFILE=throughput-performance
Environment=TUNED_IDLE_PROFILE=powersave

[Install]
WantedBy=multi-user.target
