# src/IskanderHearth/software/thermal/gpu_thermal_manager.service
# gpu_thermal_manager.service
# Iskander Hearth — GPU Thermal Throttle Manager
#
# Wraps gpu_thermal_manager.sh as a systemd service.
# Polls nvidia-smi every 5s; adjusts power limits to maintain GPU < 83°C.
#
# Install:
#   sudo cp gpu_thermal_manager.service /etc/systemd/system/
#   sudo chmod +x /opt/iskander/software/thermal/gpu_thermal_manager.sh
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now gpu_thermal_manager.service

[Unit]
Description=Iskander Hearth — GPU Thermal Throttle Manager (RTX 3060 Llama 3)
Documentation=file:///opt/iskander/software/thermal/gpu_thermal_manager.sh
After=nvidia-persistenced.service multi-user.target
Wants=nvidia-persistenced.service

[Service]
Type=simple

# Runs as root — required by nvidia-smi --power-limit
User=root

ExecStart=/opt/iskander/software/thermal/gpu_thermal_manager.sh

Restart=on-failure
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=3

# Configuration (override to tune for specific GPU TDP)
Environment=HEARTH_GPU_FLOOR_W=100
Environment=HEARTH_GPU_STEP_W=10
Environment=HEARTH_TEMP_THROTTLE=80
Environment=HEARTH_TEMP_CLAMP=85
Environment=HEARTH_TEMP_COOL=70
Environment=HEARTH_COOL_HOLD_SEC=60
Environment=HEARTH_IDLE_HOLD_SEC=300
Environment=HEARTH_POLL_SEC=5
Environment=HEARTH_OLLAMA_URL=http://127.0.0.1:11434

# Watchdog: script must ping within 60s or it gets restarted
WatchdogSec=60

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hearth.gpu_thermal

[Install]
WantedBy=multi-user.target
