# src/IskanderHearth/software/glass_box/systemd/hearth-leds.service
# hearth-leds.service
# Iskander Hearth — Glass Box NeoPixel LED State Daemon
#
# Drives WS2812B LED strip (top panel) to visualize AI operational state.
# Listens on Unix socket for state updates from all Iskander OS components.
# LED state must update within 500ms of triggering OS event.
#
# Install:
#   sudo cp hearth-leds.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now hearth-leds.service

[Unit]
Description=Iskander Hearth — Glass Box NeoPixel LED State Daemon
Documentation=file:///opt/iskander/software/glass_box/led_state_daemon.py
After=hearth-sovereignty.service
Wants=hearth-sovereignty.service

[Service]
Type=notify
NotifyAccess=main

User=hearth-ux
Group=hearth-ux

# WS2812B via rpi_ws281x requires access to /dev/mem (GPIO + DMA)
# This is handled via gpio group and the following capabilities
SupplementaryGroups=gpio

ExecStart=/opt/iskander/venv/bin/python3 \
    /opt/iskander/software/glass_box/led_state_daemon.py

ExecStop=/bin/sh -c 'kill -TERM $MAINPID'

Restart=on-failure
RestartSec=3
StartLimitIntervalSec=60
StartLimitBurst=5

# Environment
Environment=HEARTH_LED_COUNT=30
Environment=HEARTH_LED_GPIO_PIN=18
Environment=HEARTH_LED_FREQ_HZ=800000
Environment=HEARTH_LED_DMA=10
Environment=HEARTH_LED_BRIGHTNESS=80
Environment=HEARTH_LED_INVERT=0
Environment=HEARTH_LED_FRAME_HZ=60
Environment=HEARTH_LED_STATE_SOCK=/run/iskander/led_state.sock
Environment=HEARTH_SOVEREIGNTY_LOCK=/run/iskander/sovereignty.lock

RuntimeDirectory=iskander
RuntimeDirectoryMode=0750

# rpi_ws281x requires raw hardware access for PWM/DMA
# These are the minimal capabilities needed
AmbientCapabilities=CAP_SYS_RAWIO CAP_NET_ADMIN
CapabilityBoundingSet=CAP_SYS_RAWIO CAP_NET_ADMIN

NoNewPrivileges=no
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/run/iskander /dev

WatchdogSec=30

[Install]
WantedBy=multi-user.target
