# src/IskanderHearth/software/glass_box/systemd/hearth-buttons.service
# hearth-buttons.service
# Iskander Hearth — Glass Box HITL Button Handler
#
# Monitors Veto and Sign arcade buttons via GPIO.
# Sends governance signals to multi-sig daemon via Unix socket.
# Controls button illumination LEDs.
#
# Install:
#   sudo cp hearth-buttons.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now hearth-buttons.service

[Unit]
Description=Iskander Hearth — HITL Glass Box Button Handler (Veto / Sign)
Documentation=file:///opt/iskander/software/glass_box/button_handler.py
After=hearth-sovereignty.service
Wants=hearth-sovereignty.service

# Must be ready before Web3 signing daemon accepts transactions
Before=hearth-web3-signer.service

[Service]
Type=notify
NotifyAccess=main

User=hearth-ux
Group=hearth-ux

SupplementaryGroups=gpio

ExecStart=/opt/iskander/venv/bin/python3 \
    /opt/iskander/software/glass_box/button_handler.py

Restart=on-failure
RestartSec=3
StartLimitIntervalSec=30
StartLimitBurst=5

# Environment
Environment=HEARTH_GPIO_BTN_VETO=5
Environment=HEARTH_GPIO_BTN_SIGN=6
Environment=HEARTH_GPIO_LED_VETO=13
Environment=HEARTH_GPIO_LED_SIGN=19
Environment=HEARTH_SIGN_HOLD_SEC=3.0
Environment=HEARTH_EMERGENCY_HOLD_SEC=10.0
Environment=HEARTH_BTN_DEBOUNCE_SEC=0.05
Environment=HEARTH_MULTISIG_SOCK=/run/iskander/multisig.sock
Environment=HEARTH_LED_STATE_SOCK=/run/iskander/led_state.sock

RuntimeDirectory=iskander
RuntimeDirectoryMode=0750

# Hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/run/iskander

# Allow shutdown (for emergency 10s hold)
AmbientCapabilities=CAP_SYS_BOOT
CapabilityBoundingSet=CAP_SYS_BOOT

WatchdogSec=30

[Install]
WantedBy=multi-user.target
